# etc/profile.d/colorprompt.sh
# Check for interactive bash
if [ -n "${BASH_VERSION-}" ]; then
    force_color_prompt=yes
fi
