# etc/profile.d/more-ls-aliases.sh
#/bin/bash

alias ll='ls -l'
alias la='ls -A'
alias lo='ls -o'
alias lag='ls -laGg'
alias l='ls -CF'


